# platforms/nginx-rust-1.9/docker/config/supervisor/conf.d-sample/axum-api.conf
[program:axum-api]
command=/var/www/target/release/your_binary_name
directory=/var/www
user=rustuser
autostart=true
autorestart=true
stdout_logfile=/var/www/axum.log
stderr_logfile=/var/www/axum.err