# platforms/nginx-rust-1.9/docker/config/supervisor/conf.d-sample/rust-cargo-dev.conf
[program:rust_api]
directory=/var/www
command=/home/appuser/.cargo/bin/cargo run
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api_err.log