# platforms/nginx-zig/docker/config/crontab
* * * * * cd /var/www && {here-goes-worker-command} >> /dev/stdout 2>&1
